# infra/host/ib-spice@.service
# /etc/systemd/system/ib-spice@.service
# Templated SPICE websocket bridge for the in-browser console, one instance per VM
# SPICE port: `ib-spice@5900` (healer, WS 5959), `ib-spice@5910` (craft1, WS 5969),
# `ib-spice@5920` (craft2, WS 5979). The web console.html connects to ws://<host>:<WS>.
# WS port = SPICE port + 59 (see ib-spice-bridge). Unauthenticated -> LAN-only (ufw
# allow the WS ports from 10.0.0.0/24). Supersedes the single ib-spice.service.
[Unit]
Description=ib SPICE web bridge for VM SPICE %i
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jbaker
ExecStart=/home/jbaker/ib-app/ib-spice-bridge %i
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
